# symfony/m365-admin/.docker/cron/crontab
2 0 * * * root  /usr/bin/php /var/www/html/bin/console microsoft:clear-group-by-id $MICROSOFT_DEFAULT_GROUP_ID >>/tmp/m365-clean.log 2>&1
